# capsula-assistant.service
[Unit]
Description=Capsula Assistant Telegram Bot
Documentation=https://github.com/your-repo/capsula-assistant
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/Capsula_assist/capsula_assist
Environment=PATH=/Capsula_assist/capsula_assist/venv/bin
Environment=TZ=Europe/Moscow
ExecStart=/Capsula_assist/capsula_assist/venv/bin/python /Capsula_assist/capsula_assist/run.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=capsula-assistant

# Настройки безопасности
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/Capsula_assist/capsula_assist

# Ограничения ресурсов
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
